# chroot/etc/profile.d/hypernova.sh
#!/usr/bin/env bash

#
# HyperNova server management framework
#
# Dirty great environment hacks
#
# Copyright (c) 2012 TDM Ltd
#                    Luke Carrier <luke.carrier@tdm.info>
#

export PATH="$PATH:/usr/local/hypernova/bin"
